# cookbooks/lightning_node/templates/charge-lnd@.service.erb
[Unit]
Description=charge-lnd (%I)

[Service]
Type=oneshot

ExecStartPre=charge-lnd \
  --macaroon /var/charge-lnd/macaroon/charge-lnd.macaroon \
  --electrum-server localhost:50001 \
  --check \
  --config %I

ExecStartPre=date

ExecStart=charge-lnd \
  --macaroon /var/charge-lnd/macaroon/charge-lnd.macaroon \
  --electrum-server localhost:50001 \
  --config %I

StandardOutput=append:<%= @log_file %>
StandardError=append:<%= @log_file %>
