# cookbooks/lightning_node/templates/lnd.service.erb
[Unit]
Description=LND Lightning Network Daemon
<% unless @neutrino_mode %>
Wants=bitcoind.service
After=bitcoind.service <%= @after_units&.join(' ') %>
<% else %>
After=network.target <%= @after_units&.join(' ') %>
<% end %>

[Service]
# Service execution
###################
ExecStart=<%= @lnd_command %>
ExecStop=/usr/local/bin/lncli --macaroonpath /var/lnd/macaroon/admin.macaroon stop

# Process management
####################
Type=notify
Restart=always
RestartSec=30
TimeoutSec=240
LimitNOFILE=128000

# Directory creation and permissions
####################################
User=lnd

# /run/lnd
RuntimeDirectory=lnd
RuntimeDirectoryMode=0710

# Hardening measures
####################
# Provide a private /tmp and /var/tmp.
PrivateTmp=true

# Mount /usr, /boot/ and /etc read-only for the process.
ProtectSystem=full

# Disallow the process and all of its children to gain
# new privileges through execve().
NoNewPrivileges=true

# Use a new /dev namespace only populated with API pseudo devices
# such as /dev/null, /dev/zero and /dev/random.
PrivateDevices=true

# Deny the creation of writable and executable memory mappings.
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
