# cookbooks/bitcoin_node/templates/fulcrum.service.erb
[Unit]
Description=Fulcrum
PartOf=bitcoind.service
After=bitcoind.service <%= @after_units&.join(' ') %>
StartLimitBurst=2
StartLimitIntervalSec=20

[Service]
ExecStart=/usr/local/bin/Fulcrum /var/fulcrum/fulcrum.conf
KillSignal=SIGINT
User=fulcrum
Type=exec
TimeoutStopSec=300
RestartSec=30
Restart=on-failure

[Install]
WantedBy=multi-user.target
